# crontab
* * * * * root . /root/project_env.sh; bdd >> /var/log/cron.log 2>&1
